# systemd/easyfuse.socket
[Unit]
    Description=easyfuse Docker Volume Plugin socket
    PartOf=easyfuse.service

[Socket]
    ListenStream=/run/docker/plugins/easyfuse.sock
    SocketMode=0600

[Install]
    WantedBy=docker.service
